# systemd/anime-factory-yume-tick.service
[Unit]
Description=Yume factory — публикация ОДОБРЕННОЙ картинки в Telegram (Yume_Girls_AI)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ratelhhhh/anime-meme-factory-go
ExecStart=/home/ratelhhhh/anime-meme-factory-go/factory --config config.yume.json tick
