# systemd/anime-factory-tick.timer
[Unit]
Description=Запуск публикации картинки раз в час

[Timer]
OnBootSec=3min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
