# systemd/anime-factory-refill.service
[Unit]
Description=Anime meme factory — наполнение очереди постами Пикабу
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ratelhhhh/anime-meme-factory-go
ExecStart=/home/ratelhhhh/anime-meme-factory-go/factory refill
